# docker/crontab
15 * * * * flock -n /run/map_update.lock symfony app:update
