# docker/supervisor/valkey.ini
[program:valkey]
command=valkey-server
directory=/var/valkey
